# agent/secure-forwarder/secure-forwarder.service
[Unit]
Description=Secure Log Forwarder (HMAC gate for /v1/logs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root

# 작업 디렉터리 (코드 기준 위치)
WorkingDirectory=/home/last/lastagent/forwarder

# 환경 변수 파일 (.env)
EnvironmentFile=/home/last/lastagent/etc/.env

# 실행 명령 (venv 기반 Python)
ExecStart=/home/last/lastagent/venv/bin/python /home/last/lastagent/forwarder/secure-forwarder.py

Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
PrivateDevices=true
ReadOnlyPaths=/home/last/lastagent/forwarder

[Install]
WantedBy=multi-user.target
