# agent/agent.controller/agent-controller.service
[Unit]
Description=Agent Controller (pull commands from control API)
After=network-online.target secure-forwarder.service
Wants=network-online.target

[Service]
Type=simple
User=root

# 작업 디렉터리 (코드 기준 위치)
WorkingDirectory=/home/last/lastagent/agent

# 환경 변수 로딩
EnvironmentFile=/home/last/lastagent/etc/.env

# 실행 명령 (venv 기반 Python)
ExecStart=/home/last/lastagent/venv/bin/python /home/last/lastagent/agent/agent-controller.py

Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
RestrictSUIDSGID=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=/etc/secure-log-agent/remote.d
ReadWritePaths=/run
ReadWritePaths=/var/run

[Install]
WantedBy=multi-user.target
