# agent/etc/otel-agent.service
[Unit]
Description=Secure Log Agent (filelog -> OTLP/HTTP via Secure Forwarder)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# 일단 동작 우선: root 계정으로 실행
User=root
Group=root
# 나중에 하드닝하고 싶으면 아래 두 줄을 쓰고 위를 주석 처리:
#User=otel-agent
#Group=nogroup

WorkingDirectory=/etc/secure-log-agent

EnvironmentFile=/home/last/lastagent/etc/.env

PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir -p /var/lib/secure-log-agent/queue /var/lib/otelcol-contrib
ExecStartPre=/usr/bin/chown -R root:root /etc/secure-log-agent /var/lib/secure-log-agent /var/lib/otelcol-contrib

ExecStart=/usr/local/bin/otelcol-contrib \
  --config=/etc/secure-log-agent/agent.yaml

Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

ReadWritePaths=/etc/secure-log-agent
ReadWritePaths=/var/lib/otelcol-contrib
ReadWritePaths=/var/log
ReadWritePaths=/var/lib/secure-log-agent

[Install]
WantedBy=multi-user.target
