# agent/etc/otel-agent.service
[Unit]
Description=Secure Log Agent (filelog -> OTLP/HTTP via Secure Forwarder)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# otel-agent 전용 계정
User=otel-agent
Group=nogroup

WorkingDirectory=/etc/secure-log-agent

# .env 에서 INGEST_ENDPOINT / INGEST_TOKEN 등 읽어오기
EnvironmentFile=/home/last/lastagent/etc/.env

# ExecStartPre 는 root 권한으로 실행 (디렉터리/권한 보정)
PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir -p /var/lib/secure-log-agent/queue /var/lib/otelcol-contrib
ExecStartPre=/usr/bin/chown -R otel-agent:nogroup /etc/secure-log-agent /var/lib/secure-log-agent /var/lib/otelcol-contrib

# 실제 OTEL Collector 실행
ExecStart=/usr/local/bin/otelcol-contrib \
  --config=/etc/secure-log-agent/agent.yaml

Restart=always
RestartSec=5

########## Hardening ##########
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# OTEL 이 쓸 수 있는 경로만 쓰기 허용
ReadWritePaths=/etc/secure-log-agent
ReadWritePaths=/var/lib/otelcol-contrib
ReadWritePaths=/var/log
ReadWritePaths=/var/lib/secure-log-agent

[Install]
WantedBy=multi-user.target
